# ansible/roles/git-sync/templates/git-sync@.timer.j2
[Unit]
Description=Run git-sync every %i minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=%imin
Unit=git-sync.service

[Install]
WantedBy=timers.target
